# systemd/numbers-socialprove.service
[Unit]
Description=Numbers Protocol Reference Agent — SocialProve
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/opt/numbers-agents
EnvironmentFile=/opt/numbers-agents/.env
ExecStart=/opt/numbers-agents/venv/bin/python -u /opt/numbers-agents/socialprove.py
Restart=on-failure
RestartSec=15s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ref-socialprove

[Install]
WantedBy=multi-user.target
